# opt-bose-fs/etc/udev/rules.d/unlock-rndis1.rules
KERNEL=="rndis1", ACTION=="remove", SUBSYSTEM=="net", RUN+="/bin/systemctl stop unlock.service"
